# config/rc.local
#!/bin/sh -e
# rc.local for HiveMind
cd /home/pi/HiveMind-Basic
sudo python HiveMind.py
exit 0
