# config/rc.local
#!/bin/sh -e
# rc.local for HiveMind
cd /home/pi/HiveMind_Basic
sudo python HiveMind_Basic.py
exit 0
