# config/rc.local
#!/bin/sh -e
# rc.local for HiveMind
sudo python HiveMind.py
exit 0
